# .config/systemd/user/bwsync.service
[Unit]
Description="Sync Vaultwarden database"

[Service]
Type=simple
ExecStart=/usr/bin/rbw sync
