# .config/systemd/user/newsboat-reload.service
[Unit]
Description="Refresh Newsboat feed"

[Service]
Type=simple
ExecStart=/usr/bin/newsboat -x reload
